# systemd/pi-cron-scheduler.service
[Unit]
Description=Pi Cron Scheduler — Thetis Extension
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=@@EXT_DIR@@/scripts/start-scheduler.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
